# systemd/telegram-claude-bridge.service
[Unit]
Description=Telegram Claude CLI Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/projects/claudeBot
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env python3 %h/projects/claudeBot/service_entry.py --env %h/.config/telegram-claude-bridge/env
Restart=on-failure
RestartSec=3
TimeoutStopSec=20

[Install]
WantedBy=default.target
